# files/usr/lib/systemd/system/nullmailer-smtpd.socket
[Unit]
Description=Nullmailer SMTP socket

[Socket]
ListenStream=[::1]:25
Accept=yes

[Install]
WantedBy=sockets.target